# .config/macos-theme-sync/com.user.terminal-theme-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.terminal-theme-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/chris/.config/macos-theme-sync/sync.sh</string>
    </array>
    <key>WatchPaths</key>
    <array>
        <string>/Users/chris/Library/Preferences/.GlobalPreferences.plist</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
